# centos/install/commander-files/commander.service
[Unit]
Description=Commander SystemD

[Service]
ExecStart=/opt/ninja/commander

[Install]
WantedBy=multi-user.target